# etc/supervisord/nginx.ini
[program:nginx]
command=nginx.sh
stopasgroup=true
